# deploy/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock   ; (the path to the socket file)

[supervisord]
logfile=/tmp/supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10           ; (num of main logfile rotation backups;default 10)
loglevel=info                ; (log level;default info; others: debug,warn,trace)
pidfile=/tmp/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false               ; (start in foreground if true;default false)
minfds=1024                  ; (min. avail startup file descriptors;default 1024)
minprocs=200                 ; (min. avail process descriptors;default 200)
environment=
    PORT=8000,
    WEB_PASSWORD={{webpass}},
    SECRET_KEY='{{secret_key}}',
    APPLICATION_SETTINGS=production.py,
    MAIL_SERVER='{{mail_server}}',
    MAIL_PORT={{mail_port}},
    MAIL_USERNAME='{{mail_username}}',
    MAIL_PASSWORD='{{mail_password}}',
    MAIL_DEFAULT_SENDER='{{mail_default_sender}}',
    MAILER_DEBUG={{mailer_debug}},
    MAIL_DEFAULT_TO='{{mail_default_to}}',
    MAIL_DEFAULT_LIST='{{mail_default_list}}',
    DATA_ROOT='{{ data_root }}',
    ARCHIVE_PATH='{{ archive_path }}',
    RSYNC_SSH_USER='{{ rsync_ssh_user }}',
    RSYNC_HOST='{{ rsync_host }}',
    RSYNC_REMOTE_PATH='{{ rsync_remote_path }}',
    RSYNC_TO_PATH='{{ rsync_to_path }}',
    MONGO_URI='{{ mongo_db }}',
    ADMINS='{{ admins }}',
    USER_DB_FILE='{{ user_db_file }}'


[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL  for a unix socket

[program:web]
command=gunicorn app:app -b 0.0.0.0:8000 -w 2
numprocs=1
directory=/home/gliderweb/glider-mission
autostart=false
redirect_stderr=true
stdout_logfile=logs/web.log

[program:monitor_db_sync]
command=python glider_mission_db_sync.py
numprocs=1
directory=/home/gliderweb/glider-mission
stopsignal=TERM
autostart=false
redirect_stderr=true
stdout_logfile=logs/db-sync-monitor.log

