# deploy/supervisord-catalog-monitor.conf
[unix_http_server]
file=/tmp/supervisor-catalog-monitor.sock   ; (the path to the socket file)

[supervisord]
logfile=/tmp/supervisord-catalog-monitor.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10           ; (num of main logfile rotation backups;default 10)
loglevel=info                ; (log level;default info; others: debug,warn,trace)
pidfile=/tmp/supervisord-catalog-monitor.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false               ; (start in foreground if true;default false)
minfds=1024                  ; (min. avail startup file descriptors;default 1024)
minprocs=200                 ; (min. avail process descriptors;default 200)
environment=
    DATA_ROOT='{{ rsync_to_path }}',

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor-catalog-monitor.sock ; use a unix:// URL  for a unix socket

[program:catalog-monitor]
command=python glider_catalog_monitor.py
numprocs=1
directory=/home/gliderweb/glider-mission
stopsignal=TERM
autostart=false
redirect_stderr=true
stdout_logfile=logs/catalog-monitor.log

